# .claude/launchd/com.secondbrain.slackbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.secondbrain.slackbot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/viditparashar/Desktop/second-brain-starter/.venv/bin/python</string>
        <string>/Users/viditparashar/Desktop/second-brain-starter/.claude/chat/bot.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/viditparashar/Desktop/second-brain-starter</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_PROJECT_DIR</key>
        <string>/Users/viditparashar/Desktop/second-brain-starter</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/viditparashar/Desktop/second-brain-starter/.claude/data/logs/slackbot-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/viditparashar/Desktop/second-brain-starter/.claude/data/logs/slackbot-stderr.log</string>
</dict>
</plist>
